# dokku-github-webhook.conf
start on runlevel [2345]

exec sudo -u dokku /home/dokku/dokku-github-webhook/worker.rb
